# daemon/com.seyoweb.premierseyostudio.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seyoweb.premierseyostudio.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>__DAEMON_PATH__</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PREMIERSEYO_PORT</key>
    <string>53117</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/premierseyo-daemon.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/premierseyo-daemon.error.log</string>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>
</dict>
</plist>
